# foxboron-system/msmtp-update-fingerprints.service
[Unit]
Description=msmtp update fingerprints

[Service]
Type=oneshot
ExecStart=/usr/bin/msmtp-update-fingerprints
