# foxboron-system/msmtp-update-fingerprints.timer
[Unit]
Description=Syncer fingerprints from msmtp

[Timer]
Persistent=true
OnBootSec=30s
OnCalendar=hourly

[Install]
WantedBy=default.target
