# foxboron-system/timers-timer@.target
[Unit]
Description=%i timer
StopWhenUnneeded=true
